# deploy/com.daytrace.daily.plist.template
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  DayTrace daily catchup + archive export.
  Fires at 04:30 every day (after the 04:00 shifted-day boundary).
  If the Mac is asleep at fire time, launchd runs the job the next time
  the system wakes (KeepAlive default behavior for StartCalendarInterval).

  Install: run `scripts/install_launchd.sh` from the repo root. It
  substitutes __PYTHON__ and __REPO__ in this template and bootstraps
  the resulting plist into ~/Library/LaunchAgents/.

  Inspect:
    launchctl list | grep daytrace
    tail -f __REPO__/data/logs/daily.log

  Unload:
    launchctl bootout gui/$(id -u) ~/Library/LaunchAgents/com.daytrace.daily.plist
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.daytrace.daily</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>__REPO__/scripts/daytrace-daily.sh</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>4</integer>
    <key>Minute</key>
    <integer>30</integer>
  </dict>

  <key>WorkingDirectory</key>
  <string>__REPO__</string>

  <key>StandardOutPath</key>
  <string>__REPO__/data/logs/daily.log</string>
  <key>StandardErrorPath</key>
  <string>__REPO__/data/logs/daily.log</string>

  <key>RunAtLoad</key>
  <false/>
</dict>
</plist>
